# etc/sv/bamboo_cleanup/run
#!/bin/sh

exec $BAMBOO_AGENT_HOME/clearAgentSpace.sh >> /var/log/bamboo_cleanup.log 2>&1