# etc/sv/bamboo_cleanup/run
#!/bin/sh

exec $BAMBOO_AGENT_HOME/clearAgentSpace.sh 2>&1